# deployment/systemd/cloudflared.service
[Unit]
Description=Cloudflare Tunnel for Quantum Tutor
After=network-online.target nginx.service
Wants=network-online.target

[Service]
Type=simple
User=cloudflared
Group=cloudflared
ExecStart=/usr/bin/cloudflared --config /etc/cloudflared/quantumtutor-cl.yml tunnel run
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true

[Install]
WantedBy=multi-user.target
